# src/cm.conf
[program:cm]
command=/usr/local/bin/cm -h 0.0.0.0 -p 2310 -P 2443 -l /var/log/cm.log
autostart=true
autorestart=true
stderr_logfile=/var/log/cm.err.log
stdout_logfile=/var/log/cm.out.log
